# novnc-gate/conf.d/novnc-c.conf
[supervisord]
nodaemon=true
[program: MK-MUPRIV5-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65281 --vnc MK-MUPRIV5:5900
[program: MK-PRSVSHN36-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65282 --vnc MK-PRSVSHN36:5900
[program: MK-PULK15-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65283 --vnc MK-PULK15:5900
[program: MK-PULK-42-6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65284 --vnc MK-PULK-42-6:5900
[program: MK-PULK71-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65285 --vnc MK-PULK71-2:5900
[program: MK-PAYT8-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65286 --vnc MK-PAYT8:5900
[program: MK-REVOL17-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65287 --vnc MK-REVOL17:5900
[program: MK-RSTVSK24-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65288 --vnc MK-RSTVSK24:5900
[program: MK-RSNVSK18-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65289 --vnc MK-RSNVSK18:5900
[program: MK-RCHVSK13-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65290 --vnc MK-RCHVSK13:5900
[program: MK-RCHVSK4-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65291 --vnc MK-RCHVSK4:5900
[program: MK-SVER1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65292 --vnc MK-SVER1:5900
[program: RMK-SVERDLOVO1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65293 --vnc RMK-SVERDLOVO1:5900
[program: MK-SVETLA66-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65294 --vnc MK-SVETLA66:5900
[program: MK-SIREN-B17-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65295 --vnc MK-SIREN-B17:5900
[program: MK-SLAVY30-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65296 --vnc MK-SLAVY30:5900
[program: MK-SLAVYA21-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65297 --vnc MK-SLAVYA21:5900
[program: MK-SOVET36-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65298 --vnc MK-SOVET36-1:5900
[program: MK-SOLID9-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65299 --vnc MK-SOLID9:5900
[program: MK-SOFSKYA57-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65300 --vnc MK-SOFSKYA57:5900
[program: MK-SOFKOV1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65301 --vnc MK-SOFKOV1:5900
[program: MK-SRDNHTNSK25-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65302 --vnc MK-SRDNHTNSK25:5900
[program: MK-SREGAT16-4-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65303 --vnc MK-SREGAT16-4:5900
[program: MK-STRTL6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65304 --vnc MK-STRTL6:5900
[program: MK-TRFN-DOR15-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65305 --vnc MK-TRFN-DOR15:5900
[program: MK-TREF22-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65306 --vnc MK-TREF22:5900
[program: MK-UDAR30-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65307 --vnc MK-UDAR30-1:5900
[program: MK-URLSK21-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65308 --vnc MK-URLSK21:5900
[program: MK-UTOCHK6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65309 --vnc MK-UTOCHK6:5900
[program: MK-USHIN25-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65310 --vnc MK-USHIN25-1:5900
[program: MK-ABRAMOVA18-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65311 --vnc MK-ABRAMOVA18:5900
[program: MK-ABRAM21-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65312 --vnc MK-ABRAM21:5900
[program: MK-ABRAM8-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65313 --vnc MK-ABRAM8:5900
[program: MK-HUDOSH43-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65314 --vnc MK-HUDOSH43:5900
[program: MK-CHUGU1-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65315 --vnc MK-CHUGU1-2:5900
[program: MK-SH-LVRK14-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65316 --vnc MK-SH-LVRK14:5900
[program: MK-MUSHU15-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65317 --vnc MK-MUSHU15:5900
[program: MK-SHVLV18-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65318 --vnc MK-SHVLV18:5900
[program: MK-MUSHUV25-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65319 --vnc MK-MUSHUV25:5900
[program: MK-MUSHU5-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65320 --vnc MK-MUSHU5:5900
[program: MK-SHGLV9-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65321 --vnc MK-SHGLV9:5900
[program: MK-SHERB11-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65322 --vnc MK-SHERB11:5900
