# novnc-gate/conf.d/novnc-d.conf
[supervisord]
nodaemon=true
[program: RM-BUH1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65080 --vnc RM-BUH1:5900
[program: RM-BUH2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65081 --vnc RM-BUH2:5900
[program: RM-BUH5-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65082 --vnc RM-BUH5:5900
[program: RM-BUH6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65083 --vnc RM-BUH6:5900
[program: RM-BIER1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65084 --vnc RM-BIER1:5900
[program: RM-BIER2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65085 --vnc RM-BIER2:5900
[program: RM-BIER4-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65086 --vnc RM-BIER4:5900
[program: RM-BIER6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65087 --vnc RM-BIER6:5900
[program: RM-BIER8-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65088 --vnc RM-BIER8:5900
[program: RM-BIER9-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65089 --vnc RM-BIER9:5900
[program: RM-BIER7-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65090 --vnc RM-BIER7:5900
[program: RM-GROWTH1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65091 --vnc RM-GROWTH1:5900
[program: RM-LERN0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65092 --vnc RM-LERN0:5900
[program: RM-LERN1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65093 --vnc RM-LERN1:5900
[program: RM-LERN3-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65094 --vnc RM-LERN3:5900
[program: RM-LERN2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65095 --vnc RM-LERN2:5900
[program: RM-SNACKS1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65096 --vnc RM-SNACKS1:5900
[program: RM-SNACKS2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65097 --vnc RM-SNACKS2:5900
[program: RM-SNAKS4-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65098 --vnc RM-SNAKS4:5900
[program: RM-GROWTH0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65099 --vnc RM-GROWTH0:5900
[program: RM-SNACKS5-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65100 --vnc RM-SNACKS5:5900
[program: RM-SNACKS0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65101 --vnc RM-SNACKS0:5900
[program: RM-ZAKUSKI2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65102 --vnc RM-ZAKUSKI2:5900
[program: RM-ZAKUPKI1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65103 --vnc RM-ZAKUPKI1:5900
[program: RM-ZAKUPKI2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65104 --vnc RM-ZAKUPKI2:5900
[program: RM-ZAKUPKIFAEL-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65105 --vnc RM-ZAKUPKIFAEL:5900
[program: RM-MARK0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65106 --vnc RM-MARK0:5900
[program: RM-MARK1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65107 --vnc RM-MARK1:5900
[program: RM-OK0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65108 --vnc RM-OK0:5900
[program: RM-OK4-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65109 --vnc RM-OK4:5900
[program: RM-OK1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65110 --vnc RM-OK1:5900
[program: RMK-OK2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65111 --vnc RMK-OK2:5900
[program: RM-OK3-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65112 --vnc RM-OK3:5900
[program: RM-BUH0-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65113 --vnc RM-BUH0:5900
[program: RM-BUH50-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65114 --vnc RM-BUH50:5900
[program: RM--ALINA-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65115 --vnc RM--ALINA:5900
[program: RM-MARK2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65116 --vnc RM-MARK2:5900
[program: RM-MARKETING1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65117 --vnc RM-MARKETING1:5900
[program: DESKTOP-V7VL19R-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65118 --vnc DESKTOP-V7VL19R:5900
[program: OFFICEMANAGER3-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65119 --vnc OFFICEMANAGER3:5900
[program: RM-OFICCEMANAGE-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65120 --vnc RM-OFICCEMANAGE:5900
[program: SKLADALEXANDR-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65121 --vnc SKLADALEXANDR:5900
