# novnc-gate/conf.d/novnc-b.conf
[supervisord]
nodaemon=true
[program: MK-KOLON10-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65221 --vnc MK-KOLON10:5900
[program: MK-KLPNSK34-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65222 --vnc MK-KLPNSK34:5900
[program: MK-KOLP32-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65223 --vnc MK-KOLP32:5900
[program: MK-KOMEN55-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65224 --vnc MK-KOMEN55-1:5900
[program: MK-KMNDNTSK-64-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65225 --vnc MK-KMNDNTSK-64:5900
[program: MK-KMNDNTSK24-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65226 --vnc MK-KMNDNTSK24:5900
[program: MK-KMNDNTSK59-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65227 --vnc MK-KMNDNTSK59:5900
[program: MK-KOMEN67-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65228 --vnc MK-KOMEN67:5900
[program: MK-KORZU34-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65229 --vnc MK-KORZU34-1:5900
[program: MK-KOROL21-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65230 --vnc MK-KOROL21-1:5900
[program: MK-KOROLEV27-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65231 --vnc MK-KOROLEV27-1:5900
[program: MK-KOSIG27-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65232 --vnc MK-KOSIG27-1:5900
[program: MK-KRSNPUT14-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65233 --vnc MK-KRSNPUT14:5900
[program: MK-KRILK3-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65234 --vnc MK-KRILK3-1:5900
[program: MK-KUPCH15-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65235 --vnc MK-KUPCH15:5900
[program: MK-KUSHEL1-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65236 --vnc MK-KUSHEL1-2:5900
[program: MK-SH-LVRK63-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65237 --vnc MK-SH-LVRK63:5900
[program: MK-L-GLKV47-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65238 --vnc MK-L-GLKV47:5900
[program: MK-KULEN7-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65239 --vnc MK-KULEN7:5900
[program: MK-LENINS49-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65240 --vnc MK-LENINS49:5900
[program: MK-LENINS56-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65241 --vnc MK-LENINS56:5900
[program: MK-LNNSK72-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65242 --vnc MK-LNNSK72:5900
[program: MK-LNNSK77-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65243 --vnc MK-LNNSK77:5900
[program: MK-LENSOVETA-80-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65244 --vnc MK-LENSOVETA-80:5900
[program: MK-LNSVT23-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65245 --vnc MK-LNSVT23:5900
[program: MK-LIHOL14-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65246 --vnc MK-LIHOL14-2:5900
[program: MK-MRSHK10-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65247 --vnc MK-MRSHK10:5900
[program: MK-MARSH12-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65248 --vnc MK-MARSH12:5900
[program: MK-M-ZKHRV21-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65249 --vnc MK-M-ZKHRV21:5900
[program: MK-KZKV70-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65250 --vnc MK-KZKV70:5900
[program: MK-MUMEND10-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65251 --vnc MK-MUMEND10:5900
[program: MK-MRKRV1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65252 --vnc MK-MRKRV1:5900
[program: MK-METAL111-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65253 --vnc MK-METAL111:5900
[program: MK-MORSKAYA35-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65254 --vnc MK-MORSKAYA35:5900
[program: MK-MURDOR-16-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65255 --vnc MK-MURDOR-16-2:5900
[program: MK-MURDOR61-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65256 --vnc MK-MURDOR61:5900
[program: MK-NALICH40-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65257 --vnc MK-NALICH40:5900
[program: MK-NAROD45-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65258 --vnc MK-NAROD45:5900
[program: MK-NRND68-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65259 --vnc MK-NRND68:5900
[program: MK-N-OPLCHN10-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65260 --vnc MK-N-OPLCHN10:5900
[program: MK-NAST35-1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65261 --vnc MK-NAST35-1:5900
[program: MK-NAST8-2-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65262 --vnc MK-NAST8-2:5900
[program: MK-NAUKI24-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65263 --vnc MK-NAUKI24:5900
[program: MK-NEVSKAYA7-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65264 --vnc MK-NEVSKAYA7:5900
[program: RM-NVTRV10-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65265 --vnc RM-NVTRV10:5900
[program: MK-MUNOV7-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65266 --vnc MK-MUNOV7:5900
[program: MK-NVGRDSK6-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65267 --vnc MK-NVGRDSK6:5900
[program: MK-NVZMLVSK13-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65268 --vnc MK-NVZMLVSK13:5900
[program: MK-NOVOCH40-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65269 --vnc MK-NOVOCH40:5900
[program: MK-KUOBL3-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65270 --vnc MK-KUOBL3:5900
[program: MK-OKTNAB118-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65271 --vnc MK-OKTNAB118:5900
[program: MK-OPTIK37-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65272 --vnc MK-OPTIK37:5900
[program: MK-PRSHTN61-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65273 --vnc MK-PRSHTN61:5900
[program: MK-PARF1-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65274 --vnc MK-PARF1:5900
[program: MK-PZZHN14-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65275 --vnc MK-PZZHN14:5900
[program: MK-PTRGFSK17-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65276 --vnc MK-PTRGFSK17:5900
[program: MK-PLNRN87-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65277 --vnc MK-PLNRN87:5900
[program: MK-PLESE16-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65278 --vnc MK-PLESE16:5900
[program: MK-KUPRZ14-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65279 --vnc MK-KUPRZ14:5900
[program: MK-KUPRZ5-vnc]
command=/usr/share/novnc/utils/launch.sh --listen 65280 --vnc MK-KUPRZ5:5900
